# scripts/jupyterasservice/jupyter.service
[Unit]
After=network.service

[Service]
ExecStart=jupyter notebook
User=geot

[Install]
WantedBy=default.target
